# .config/systemd/user/wallpaper-changer.service
# ~/.config/systemd/user/wallpaper-changer.service
[Unit]
Description=Hyprland Random Wallpaper Changer
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/hypr/scripts/random-wallpaper.sh daemon
Restart=on-failure
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target

# Run this service:
#
# 重新加载 systemd
# systemctl --user daemon-reload
#
# 启用服务（可选）
# systemctl --user enable wallpaper-changer.service
